# $HOME/.config/autostart/eww.desktop
[Desktop Entry]
Exec=/usr/local/bin/eww open leftbar
TryExec=/usr/local/bin/eww
Name=eww
Comment=Render a Linux desktop bar
Type=Application
