# $HOME/.config/autostart/stalonetray.desktop
[Desktop Entry]
Exec=/usr/bin/stalonetray
TryExec=/usr/bin/stalonetray
Name=stalonetray
Comment=Stand-alone system tray
Type=Application
