# $HOME/.config/autostart/xscreensaver.desktop
# BEGIN ANSIBLE MANAGED BLOCK
[Desktop Entry]
# Managed by the desktop role (tasks/idle.yml). The daemon reads ~/.xscreensaver at
# startup and reloads it whenever it changes.
Exec=/usr/bin/xscreensaver
TryExec=/usr/bin/xscreensaver
Name=xscreensaver
Comment=Blank, lock and power down an idle X session
Type=Application
# END ANSIBLE MANAGED BLOCK
